# ConfluenceRag/cron/crontab
* * * * * /usr/local/bin/python /update_chroma_db.py >> /var/log/cronjob.log 2>&1
