# ConfluenceRag/cron/crontab
* * * * * /usr/local/bin/python /conftest.py >> /var/log/cronjob.log 2>&1
